# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/var/www/mu-cecsc/server.py