# supervisord.conf
[supervisord]
nodaemon=false

[program:server]
command=python3 /var/www/mu-cec/server.py